# board/m5stack/overlay/etc/rc.local
#!/bin/sh
[ "$1" = "stop" ] && exit 0



















exit 0
