# board/m5stack/overlay/etc/rc.local
#!/bin/sh
/usr/local/m5stack/startup-script.sh

/soc/scripts/startup-script.sh

bash -c ". /etc/profile ; /soc/scripts/usb-adb.sh start"

chown root:telnetd /usr/lib/telnetlogin
chmod 4754 /usr/lib/telnetlogin
/usr/local/m5stack/bin/ax_usb_adb_event.sh >> /dev/null 2>&1 &
