# tools/overlay_ubuntu22_04/etc/profile.d/locale_env.sh
export LC_ALL=C.utf8
export LANGUAGE=C.utf8
export LANG=C.utf8
