# tools/overlay_ubuntu22_04/etc/profile.d/m5stack.sh
#!/bin/sh
. /usr/local/m5stack/bash_env.sh