# tools/overlay_ubuntu22_04/etc/network/if-pre-up.d/set-mac-eth0
#!/bin/sh
[ "$IFACE" = "eth0" ] || exit 0
/usr/local/m5stack/bin/S39eth_mac set_mac eth0
exit 0
