# board/m5stack/overlay/etc/profile.d/axenv.sh
#!/bin/sh
export PATH=/opt/bin:/opt/usr/bin:/opt/scripts:/soc/bin:/soc/scripts:/usr/local/bin:$PATH
export LD_LIBRARY_PATH=/opt/lib:/opt/usr/lib:/soc/lib:/usr/local/lib:/usr/lib