# tools/overlay_ubuntu22_04/etc/rc.local
#!/bin/sh
/usr/local/m5stack/startup-script.sh
/soc/scripts/startup-script.sh



chown root:telnetd /usr/lib/telnetlogin
chmod 4754 /usr/lib/telnetlogin

