# tools/overlay_ubuntu22_04/etc/rc.local
#!/bin/bash
bash /soc/scripts/auto_load_all_drv.sh
bash /usr/local/m5stack/startup-script.sh
bash /soc/scripts/startup-script.sh

bash -c ". /etc/profile ; /soc/scripts/usb-adb.sh start"
