# tools/overlay_ubuntu22_04_LITE/etc/rc.local
#!/bin/bash
/usr/local/m5stack/bin/S03espc6wifi start
bash /soc/scripts/startup-script.sh
bash /usr/local/m5stack/startup-script.sh

bash -c ". /etc/profile ; /soc/scripts/usb-adb.sh start"
exit 0